# route-ap-clients.service
[Unit]
Description=Route AP clients via Tailscale exit node
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/route-ap-clients.sh

[Install]
WantedBy=multi-user.target
