# wormhole-localwan.service
[Unit]
Description=Keep wormhole's own traffic on the local router's internet connection
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/localwan.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
